# usbip/rootfs/etc/services.d/usbip/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: USBIP Mounter
# Start the usbipd daemon service
# ==============================================================================

# Wait for usbipd to become available
s6-svwait -u -t 5000 /run/service/usbipd

# Loop attachment code every 5 minutes
bashio::log.info "Attaching usbip devices"
while true; do
    bashio::log.debug "Running mount_devices script..."
    bashio::log.debug "$(/usr/local/bin/mount_devices 2>&1)"
done
